# samples/_root_/etc/rc.local
#!/bin/bash

# starts audio system
runuser -l predic -c "/home/predic/pre.di.c/.venv/bin/python /home/predic/pre.di.c/startaudio.py"

# loads irexec remote command launcher
runuser -l predic -c "/usr/bin/irexec -d"
exit 0
