# samples/_root_/etc/rc.local
#!/bin/bash

# starts audio system
#su -l predic -c "/usr/bin/python3 /home/predic/pre.di.c/startaudio.py"
runuser -l predic -c "/usr/bin/python3 /home/predic/pre.di.c/startaudio.py"

# loads irexec remote command launcher
#su -l predic -c "/usr/bin/irexec -d"
runuser -l predic -c "/usr/bin/irexec -d"
exit 0
